# deploy/backup_sqlite.cron
0 2 * * * /opt/coaching_automation/deploy/backup_sqlite.sh >> /var/log/coaching_automation_backup.log 2>&1
